# scripts/systemd/noc-probe-runner.timer
[Unit]
Description=Run OKX Agentic Bot NOC Probe every 5 minutes

[Timer]
OnCalendar=*:0/5
Persistent=true
Unit=noc-probe-runner.service

[Install]
WantedBy=timers.target
